# docs/examples/supervisord/captiveportal.ini
[program:captiveportal]
command={{global_virtualenv_path}}/bin/gunicorn portal --workers 3 --log-level DEBUG --bind 127.0.0.1:8080
directory=/var/opt/captiveportal/portal
pythonpath={{global_virtualenv_path}}/lib/python2.7/site-packages
user={{captiveportal_username}}
group=virtualenv
umask=002
autostart=true
autorestart=true
stdout_logfile=/var/log/captiveportal_gunicorn.log
loglevel=debug
redirect_stderr=true
stopsignal=QUIT
environment=PATH={{global_virtualenv_path}}/bin:/bin:/sbin:/usr/bin:/usr/sbin
